# airootfs/etc/systemd/system/make-liveuser.service
[Unit]
Description=Create live user
After=network.target
Before=nextcloud-web.service

[Service]
ExecStart=/etc/mkliveuser.sh
